# crontab
0 0 2 * * * * cd /code && python scripts/backup_db.py # backups the database every day at 2am